# mqtt2sms-addon/rootfs/etc/services.d/mqtt2sms/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the MQTT2SMS service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here
#!/usr/bin/with-contenv bashio

# go to script directory
cd /usr/src/mqtt2sms

# extract config
cp /data/options.json config.json
chmod g+rw,o+r config.json

# extract mappings
jq '.COMMANDS2SMS_MAPPING' config.json > mappings/commands2sms/default.json
jq '.CUSTOM_SMS2STATES_MAPPING' config.json > mappings/sms2states/custom.json

# Get some gammu configuration from options
GAMMU_DEVICE=$(bashio::config 'GAMMU_DEVICE') GAMMU_PIN=$(bashio::config 'GAMMU_PIN') envsubst < /etc/gammu-smsdrc.base > /etc/gammu-smsdrc

# Start gammu daemon
start-stop-daemon --start --quiet --pidfile /var/run/gammu-smsd.pid \
		  --exec /usr/bin/gammu-smsd -- --daemon --user root --group root \
                  --pid /var/run/gammu-smsd.pid

# Execute gateways + log output from gammu
node mqtt2gammu.js & node app.js & tail -f /var/log/gammu
